# Galway-Bay/crontab
*/5 * * * * /usr/local/bin/python /root/galway.py 
