# eBird/crontab
37 15 * * * /usr/local/bin/python /root/main.py 
